# DOCKER/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn3 app:app --workers 1 --bind unix:gunicorn.sock -m 777 --log-file=gunicorn_errors.log
directory=/app/app